# debian/kernelstub.postinst
#!/bin/bash

systemctl enable kernelstub-hide-menu
